# root/etc/services.d/dnscrypt-proxy/run
#!/usr/bin/with-contenv bash
exec \
s6-setuidgid dnsx /usr/bin/dnscrypt-proxy -config /dnscrypt-proxy.toml
